# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:app-web-1]
command=gunicorn portfolio
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/app/web-1.log
stderr_logfile=/var/log/app/web-1.error.log
user=root
directory=/root/portfolio
environment=PORT="80"

[program:app-redis-1]
command=redis-server
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/app/redis-1.log
stderr_logfile=/var/log/app/redis-1.error.log
user=root
directory=/root/portfolio

[group:app]
programs=app-web-1, app-redis-1
